# openagent.service
[Unit]
Description=OpenAgent AI Assistant
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/a0/usr/workdir/openagent
Environment="NODE_ENV=production"
Environment="PATH=/root/.nvm/versions/node/v22.22.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/node src/cli.js --no-interactive
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
